# delinq.timer
[Unit]
Description=timer for delinq script 
 
[Timer]
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target

